# rootfs/etc/systemd/system/systemd-poweroff.timer
[Unit]
Description=Automatic Power-Off

[Timer]
OnCalendar=23:30
Persistent=false

[Install]
WantedBy=basic.target
